# data/systemd/redshift-gtk.service.in
[Unit]
Description=Redshift display colour temperature adjustment (GUI)
Documentation=http://jonls.dk/redshift/

[Service]
ExecStart=@bindir@/redshift-gtk
Restart=on-failure
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
